# upstart/irc-notifier.conf
# /etc/init/irc-notifier.conf

description "irc-notifier server"
author      "Christopher Sexton"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 99 5

chdir /var/www/irc-notifier

script
    exec IRC_PASSWORD="secret" IRC_ROOM="#notifier" IRC_SERVER="localhost" /usr/local/bin/bundle exec rackup -p 9001 config.ru >> /var/log/irc-notifier.log 2>&1
end script
